# deploy/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:fastapi]
command=uvicorn app.main:app --host 0.0.0.0 --port 9900
directory=/app
user=root
autostart=true
autorestart=true
stdout_logfile=/app/logs/api.log
stderr_logfile=/app/logs/api.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=
    LOGURU_LEVEL="INFO"

[program:mcp-cls]
command=python mcp_servers/cls_server.py
directory=/app
user=root
autostart=true
autorestart=true
stdout_logfile=/app/logs/cls.log
stderr_logfile=/app/logs/cls.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:mcp-monitor]
command=python mcp_servers/monitor_server.py
directory=/app
user=root
autostart=true
autorestart=true
stdout_logfile=/app/logs/monitor.log
stderr_logfile=/app/logs/monitor.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
